# rootfs/etc/services.d/cron/run
#!/usr/bin/with-contenv sh

crontab -l
cron -f
